# services/voxl-tflite-server.service.d/partof-camera.conf
[Unit]
# Restart tflite with the camera server so disparity recovers after camera death.
PartOf=voxl-camera-server.service
After=voxl-camera-server.service
