# services/voxl-camera-server.service.d/wants-depth.conf
[Unit]
# Bring disparity (tflite) and metric depth back up whenever camera-server starts.
Wants=voxl-tflite-server.service mono_depth_rescaler.service
